# com.claude.dictation-auto-submit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.dictation-auto-submit</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/allenhu/.claude/hooks/dictation-auto-submit</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/dictation-auto-submit.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/dictation-auto-submit.log</string>
</dict>
</plist>
